# roles/services/thermald/templates/thermald.service.d/ignore_cpuid_check.conf.j2
## {{ ansible_managed }}

[Service]
ExecStart=
ExecStart=/usr/bin/thermald --no-daemon --dbus-enable --ignore-cpuid-check
